# whatsapp-agent.service
[Unit]
Description=WhatsApp Dev Agent (hardened)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wa-agent
Group=wa-agent
WorkingDirectory=/opt/whatsapp-agent
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=10

# Variables d'environnement
EnvironmentFile=/etc/whatsapp-agent.env

# === Hardening systemd ===
# NoNewPrivileges + capabilities vides + ProtectSystem couvrent l'essentiel.
# On NE met PAS SystemCallFilter ni MemoryDenyWriteExecute ni RestrictRealtime
# parce que Claude Code (Bun-compiled) appelle sched_setscheduler au boot
# et fait du JIT mémoire (V8) → SIGSYS / segfault sinon.

NoNewPrivileges=true

# Filesystem
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/whatsapp-agent
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible

# Network — uniquement IPv4/IPv6 (pas de raw sockets)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Capabilities — aucune privilégiée
CapabilityBoundingSet=
AmbientCapabilities=

# Resource limits
LimitNOFILE=4096
LimitNPROC=512
MemoryMax=2G
TasksMax=200

# Defense-in-depth (ces options ne cassent pas Claude Code)
LockPersonality=true
RestrictNamespaces=true
RestrictSUIDSGID=true
RemoveIPC=true

# Auto-restart limits (anti boucle de crash)
StartLimitIntervalSec=600
StartLimitBurst=5

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=whatsapp-agent
Environment=HOME=/opt/whatsapp-agent

[Install]
WantedBy=multi-user.target
